# sipcallnotify/DEBIAN/postinst
#!/bin/sh
update-rc.d -f sipcallnotify defaults